# dockerimages/autoscaler/supervisord.conf
[supervisord]
nodaemon=true

[program:autoscaler]
command=python3 autoscaler.py

[program:response_visualizer]
command=python3 response_visualizer.py

[program:analyzer]
command=python3 analyzer.py